# scripts/sendDueRemindersTest.cron
* * * * * cd /home/teju/crm && /home/teju/.asdf/shims/node scripts/sendDueReminders.cjs >> /home/teju/crm/cron.log 2>&1
